# supervisor.conf
[program:imagr]
user=ubuntu
stdout_logfile=/home/ubuntu/sv_log
command=/usr/bin/python /home/ubuntu/imagr/manage.py runserver --configuration=Prod
directory=/home/ubuntu/imagr
autorestart=true